# src/cron.d/default
*/1 * * * * root /usr/local/bin/php /www/study/yii daemon/test
